# cron.d/zxtm
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin

# Produce graphs every 5 minutes
*/5 * * * *	root zxtm-graph &> /var/log/zxtm/graph.log

# Monitor for problems every 15 minutes
*/15  * * * *	root zxtm-monitor &> /var/log/zxtm/monitor.log
